# install/systemd/user/rsynapse-notifications.service.in
[Unit]
Description=Rsynapse Notifications
After=graphical-session.target

[Service]
Type=simple
ExecStart=@LOCAL_BIN@/rsynapse-notifications
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target

